# .config/brillo_perm/changePerm.service
[Unit]

Description=make any service access brillo
[Service]

ExecStart=/usr/local/bin/changePerm.sh

[Install]
WantedBy=multi-user.target
